# build/pkg-scripts/postinstall
#!/bin/sh

# Find the .app bundle in the installation directory
APP_BUNDLE=""
for app in "$2"/*.app; do
    if [ -d "$app" ]; then
        APP_BUNDLE="$app"
        break
    fi
done

# Fallback to CrowVPN.app if not found (for development or direct install)
if [ -z "$APP_BUNDLE" ] || [ ! -d "$APP_BUNDLE" ]; then
    APP_BUNDLE="$2/CrowVPN.app"
fi

# If still not found, try Sparkle.app (legacy support)
if [ ! -d "$APP_BUNDLE" ]; then
    APP_BUNDLE="$2/Sparkle.app"
fi

if [ ! -d "$APP_BUNDLE" ]; then
    echo "Warning: Could not find app bundle in $2" >&2
    exit 0
fi

echo "Found app bundle: $APP_BUNDLE" >&2

# Set permissions for mihomo binaries
# First, remove quarantine attributes if present
xattr -d com.apple.quarantine "$APP_BUNDLE/Contents/Resources/sidecar/mihomo" 2>/dev/null || true
xattr -d com.apple.quarantine "$APP_BUNDLE/Contents/Resources/sidecar/mihomo-alpha" 2>/dev/null || true

# Try to set ownership and SUID bit
# Note: chown may fail due to SIP, but chmod +s should work in /Applications
MIHOMO_PATH="$APP_BUNDLE/Contents/Resources/sidecar/mihomo"
MIHOMO_ALPHA_PATH="$APP_BUNDLE/Contents/Resources/sidecar/mihomo-alpha"

if [ -f "$MIHOMO_PATH" ]; then
    echo "Setting permissions for mihomo..." >&2
    # Try chown first (may fail, but that's OK)
    chown root:admin "$MIHOMO_PATH" 2>/dev/null || chown root:wheel "$MIHOMO_PATH" 2>/dev/null || true
    # Set SUID bit (this is the critical part)
    chmod +sx "$MIHOMO_PATH" 2>/dev/null || {
        echo "Warning: Failed to set SUID bit for mihomo. This may be due to SIP restrictions." >&2
        # Try without SUID, just make it executable
        chmod +x "$MIHOMO_PATH" 2>/dev/null || true
    }
    # Verify permissions
    if [ -x "$MIHOMO_PATH" ]; then
        PERMS=$(ls -l "$MIHOMO_PATH" | awk '{print $1}')
        if echo "$PERMS" | grep -q 's'; then
            echo "SUID bit set successfully for mihomo" >&2
        else
            echo "Warning: SUID bit not set for mihomo (permissions: $PERMS)" >&2
        fi
    fi
fi

if [ -f "$MIHOMO_ALPHA_PATH" ]; then
    echo "Setting permissions for mihomo-alpha..." >&2
    chown root:admin "$MIHOMO_ALPHA_PATH" 2>/dev/null || chown root:wheel "$MIHOMO_ALPHA_PATH" 2>/dev/null || true
    chmod +sx "$MIHOMO_ALPHA_PATH" 2>/dev/null || {
        echo "Warning: Failed to set SUID bit for mihomo-alpha. This may be due to SIP restrictions." >&2
        chmod +x "$MIHOMO_ALPHA_PATH" 2>/dev/null || true
    }
    if [ -x "$MIHOMO_ALPHA_PATH" ]; then
        PERMS=$(ls -l "$MIHOMO_ALPHA_PATH" | awk '{print $1}')
        if echo "$PERMS" | grep -q 's'; then
            echo "SUID bit set successfully for mihomo-alpha" >&2
        else
            echo "Warning: SUID bit not set for mihomo-alpha (permissions: $PERMS)" >&2
        fi
    fi
fi

NEW_SYSPROXY="$APP_BUNDLE/Contents/Resources/files/sysproxy"
OLD_SYSPROXY="/Library/PrivilegedHelperTools/sparkle.helper"

if [ ! -f "$NEW_SYSPROXY" ]; then
    echo "Warning: Helper file not found at $NEW_SYSPROXY" >&2
    exit 0
fi

echo "Found helper file: $NEW_SYSPROXY" >&2

NEED_UPDATE=false

if [ ! -f "$OLD_SYSPROXY" ]; then
    NEED_UPDATE=true
else
    if ! cmp -s "$NEW_SYSPROXY" "$OLD_SYSPROXY"; then
        NEED_UPDATE=true
    fi
fi

if [ "$NEED_UPDATE" = true ]; then
    echo "Installing Helper..." >&2
    mkdir -p /Library/PrivilegedHelperTools
    cp "$NEW_SYSPROXY" /Library/PrivilegedHelperTools/sparkle.helper
    chown root:wheel /Library/PrivilegedHelperTools/sparkle.helper
    chmod 544 /Library/PrivilegedHelperTools/sparkle.helper
    echo "Helper installed successfully" >&2

    cat << EOF > /Library/LaunchDaemons/sparkle.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>sparkle.helper</string>
        <key>MachServices</key>
        <dict>
                <key>sparkle.helper</key>
                <true/>
        </dict>
        <key>KeepAlive</key>
        <true/>
        <key>Program</key>
        <string>/Library/PrivilegedHelperTools/sparkle.helper</string>
        <key>ProgramArguments</key>
        <array>
                <string>/Library/PrivilegedHelperTools/sparkle.helper</string>
                <string>server</string>
        </array>
        <key>StandardErrorPath</key>
        <string>/tmp/sparkle.helper.log</string>
        <key>StandardOutPath</key>
        <string>/tmp/sparkle.helper.log</string>
    </dict>
</plist>
EOF

    chown root:wheel /Library/LaunchDaemons/sparkle.helper.plist
    chmod 644 /Library/LaunchDaemons/sparkle.helper.plist
    launchctl load /Library/LaunchDaemons/sparkle.helper.plist 2>/dev/null || true
    launchctl start sparkle.helper 2>/dev/null || true
    echo "Helper daemon started" >&2
else
    echo "Helper already up to date" >&2
fi

# Try to open the newly installed app after installation
# Wait a bit for installation to complete, then open the app
APP_PATH="/Applications/CrowVPN.app"
if [ ! -d "$APP_PATH" ]; then
    APP_PATH="/Applications/Sparkle.app"
fi

if [ -d "$APP_PATH" ]; then
    echo "Opening newly installed app: $APP_PATH" >&2
    # Wait 3 seconds for installation to fully complete, then open the app
    # Use background process to avoid blocking installation
    (sleep 3 && open "$APP_PATH" 2>/dev/null) &
fi

exit 0
